# fs/etc/profile.d/oge.sh
#!/bin/bash
export SGE_ROOT=/LUSTRE/apps/oge
source /LUSTRE/apps/oge/default/common/settings.sh
